# pkg/scripts/postinst.sh
#!/bin/sh
# Load the font_cache launchd.plist
/bin/launchctl load -w /System/Library/LaunchDaemons/org.x.font_cache.plist

exit 0



